# agent/postnow-print-agent.service
[Unit]
Description=PostNow Linux Print Agent
After=network-online.target cups.service
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=/etc/postnow-print-agent.env
ExecStart=/usr/bin/python3 /opt/postnow/linux-print-agent.py
Restart=on-failure
RestartSec=5
User=postnow-agent
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
